# quadlets/dagster-postgres.container
[Unit]
Description=PostgreSQL Database for Dagster
After=dagster-pod.service
Requires=dagster-pod.service

[Container]
Image=docker.io/library/postgres:15
ContainerName=dagster-postgres
Pod=dagster.pod
Volume=/var/lib/dagster/postgres:/var/lib/postgresql/data
Environment=POSTGRES_USER=dagster
Environment=POSTGRES_PASSWORD=dagster
Environment=POSTGRES_DB=dagster
Environment=PGDATA=/var/lib/postgresql/data/pgdata
HealthCmd=pg_isready -U dagster
HealthInterval=10s
HealthTimeout=5s
HealthRetries=5

[Service]
Restart=always
TimeoutStartSec=60

[Install]
WantedBy=default.target
