# quadlets/dnsmasq.container
[Unit]
Description=DNSmasq Nameserver for liminati.internal
After=network-online.target
Wants=network-online.target

[Container]
Image=docker.io/jpillora/dnsmasq:latest
ContainerName=dnsmasq
Volume=/var/lib/dnsmasq-llm/dnsmasq.conf:/etc/dnsmasq.conf:ro
Environment=HTTP_USER=admin
Environment=HTTP_PASS=admin
Environment=HTTP_PORT=5380
# Use host network to bind to privileged port 53
# Runs as system service (rootful podman) for DNS server
Network=host

[Service]
Restart=always
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
