# quadlets/grafana.container
[Unit]
Description=Grafana Dashboard
After=network-online.target prometheus.service
Wants=network-online.target
Requires=prometheus.service

[Container]
Image=docker.io/grafana/grafana:latest
ContainerName=grafana
Volume=%h/model-serving/config/grafana/data:/var/lib/grafana
# Use host network to access Prometheus on localhost:9090
Network=host
# Run with UID/GID mapping for rootless podman
PodmanArgs=--userns=keep-id
Environment=GF_SECURITY_ADMIN_PASSWORD=admin
Environment=GF_USERS_ALLOW_SIGN_UP=false
Environment=GF_SERVER_HTTP_PORT=3000
Environment=GF_SERVER_DOMAIN=grafana.liminati.internal
Environment=GF_SERVER_ROOT_URL=https://grafana.liminati.internal

[Service]
Restart=always
TimeoutStartSec=60

[Install]
WantedBy=default.target
