# quadlets/nginx-proxy.container
[Unit]
Description=Nginx Reverse Proxy with SSL Termination
After=network-online.target vllm-qwen.service open-webui.service
Wants=network-online.target

[Container]
Image=docker.io/library/nginx:alpine
ContainerName=nginx-proxy
PublishPort=443:443
PublishPort=80:80
Volume=/var/lib/vllm/nginx/nginx.conf:/etc/nginx/nginx.conf:ro
Volume=/var/lib/vllm/nginx/conf.d:/etc/nginx/conf.d:ro
Volume=/var/lib/vllm/nginx/ssl:/etc/nginx/ssl:ro
PodmanArgs=--network=host

[Service]
Restart=always
TimeoutStartSec=60

[Install]
WantedBy=default.target
