# quadlets/open-webui.container
[Unit]
Description=Open WebUI for vLLM
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/open-webui/open-webui:main
ContainerName=open-webui
PublishPort=3000:8080
Volume=/var/lib/webui/data:/app/backend/data
Environment=OPENAI_API_BASE_URL=http://localhost:8000/v1
Network=host

[Service]
Restart=always
TimeoutStartSec=300

[Install]
WantedBy=default.target
