# quadlets/opik-minio.container
[Unit]
Description=Opik MinIO Object Storage
After=network-online.target opik-network.service
Wants=network-online.target

[Container]
Image=docker.io/minio/minio:RELEASE.2025-03-12T18-04-18Z
ContainerName=minio
HostName=minio
Network=opik-network.network
PublishPort=9000:9000
PublishPort=9090:9090
Volume=/var/lib/opik/minio:/data:Z
Environment=MINIO_ROOT_USER=THAAIOSFODNN7EXAMPLE
Environment=MINIO_ROOT_PASSWORD=LESlrXUtnFEMI/K7MDENG/bPxRfiCYEXAMPLEKEY
Exec=server --console-address ":9090" /data

[Service]
Restart=always
TimeoutStartSec=300

[Install]
WantedBy=default.target
