# quadlets/opik-minio.container
[Unit]
Description=Opik MinIO Object Storage
After=network-online.target
Wants=network-online.target

[Container]
Image=docker.io/minio/minio:RELEASE.2025-03-12T18-04-18Z
ContainerName=opik-minio
Pod=opik.pod
Volume=/var/lib/opik/minio:/data:Z
Environment=MINIO_ROOT_USER=THAAIOSFODNN7EXAMPLE
Environment=MINIO_ROOT_PASSWORD=LESlrXUtnFEMI/K7MDENG/bPxRfiCYEXAMPLEKEY
Exec=server --console-address ":9090" /data
HealthCmd=curl -f http://localhost:9000/minio/health/live
HealthInterval=1s
HealthTimeout=10s
HealthRetries=30
HealthStartPeriod=5s

[Service]
Restart=always
TimeoutStartSec=300

[Install]
WantedBy=default.target
